# etc/supervisord.d/kernelscope.ini
[program:kernelscope]
command=python KernelscopeLoggerService.py --sqlite /var/lib/kernelscope.db 8081
directory=/opt/kernelscope/src/
redirect_stderr=true
autostart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
