# etc/supervisord.d/kernelscope-offcputime.ini
[program:kernelscope-offcputime]
command=python offcputime.py --logger 'http://localhost:8081' --threshold 5000
directory=/opt/kernelscope/src/
redirect_stderr=true
autostart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
