# etc/supervisord.d/kernelscope-visualiser.ini
[program:kernelscope-visualiser]
command=python KernelscopeService.py --sqlite /var/lib/kernelscope.db 8080
directory=/opt/kernelscope/src/
redirect_stderr=true
autostart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
